# layers/meta-plink/recipes-bsp/l4t-usb-device-mode/l4t-usb-device-mode/99-usb-network.rules
SUBSYSTEM=="net", ACTION=="add", NAME=="usb0", RUN+="/bin/sh -c 'echo device > /sys/class/usb_role/usb2-0-role-switch/role'"