# layers/meta-mli/recipes-core/fake-hwclock/files/fake-hwclock-tick.service
[Unit]
Description=Save hwclock by timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fake-hwclock tick > /dev/null
